# com.benlab.lumideck-midi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.benlab.lumideck-midi</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.11/bin/python3</string>
    <string>REPLACE_WITH_ABSOLUTE_PATH/lumideck_midi.py</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardErrorPath</key><string>/tmp/lumideck-midi.log</string>
  <key>StandardOutPath</key><string>/tmp/lumideck-midi.log</string>
</dict>
</plist>
